# ansible-junos/dhcpd-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:dhcpd]
command=/usr/sbin/dhcpd -f -cf /etc/dhcp/dhcpd.conf --no-pid eth0

[program:vsftpd]
command=/usr/sbin/vsftpd
